# ops/denue-matview-refresh.service
# Snapshot of /etc/systemd/system/denue-matview-refresh.service.
#
# Install on a fresh VPS:
#   sudo cp ops/denue-matview-refresh.{service,timer} /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now denue-matview-refresh.timer
#
# This unit wraps scripts/refresh-matviews.sh (REFRESHes the 4 analytics MVs).
# Triggered by denue-matview-refresh.timer; can also be run manually:
#   sudo systemctl start denue-matview-refresh.service
#
# Why a timer: loaders are run manually, but the handlers fall back to live
# aggregation when an MV is missing entirely. They have NO way to detect
# "MV exists but is stale." A daily REFRESH defends against forgotten
# refreshes after a load, capping staleness at 24h.
#
# Logs: `journalctl -u denue-matview-refresh -f`
# Status: `systemctl status denue-matview-refresh.timer`

[Unit]
Description=DENUE analyzer mat-view refresh (analytics endpoints)
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/root/claude/projects/data-intelligence/denue-data-analysis
ExecStart=/root/claude/projects/data-intelligence/denue-data-analysis/scripts/refresh-matviews.sh
# Audit C1 round-1 closure 2026-05-10: bumped 300→1800. Live run on
# 2026-05-09 took 194s; 2026-05-10 timed out at 300s (mv_sector_grade_matrix
# alone took 294s — right on the edge). 1800s = 30min headroom for ~8×
# growth without re-tuning.
TimeoutStartSec=1800
StandardOutput=journal
StandardError=journal
SyslogIdentifier=denue-matview-refresh

Environment=TZ=America/Mexico_City
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

[Install]
WantedBy=multi-user.target
